# tinybrain/tinybrain.service
[Unit]
Description=BUGOUT tinybrain

[Service]
ExecStart=/home/youruser/.cargo/bin/tinybrain

[Install]
WantedBy=multi-user.target
~
