# admin/config-redis.service
[Unit]
Description=Configure kernel for redis

[Service]
ExecStart=/usr/bin/sh /home/core/BUGOUT/admin/config-redis.sh

[Install]
WantedBy=multi-user.target
